# skywriter-draw.service
# Skywriter exhibition app systemd user service: start on boot + auto-restart on crash (watchdog).
#
# Install (run once on the Pi):
#   mkdir -p ~/.config/systemd/user
#   cp ~/skywriter-draw.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now skywriter-draw
#
# Common commands:
#   systemctl --user status skywriter-draw     # check status
#   journalctl --user -u skywriter-draw -f     # follow logs
#   systemctl --user stop skywriter-draw       # stop (for debugging; ESC exit gets restarted automatically, stop first)
#   systemctl --user disable skywriter-draw    # disable start on boot
#
# Prerequisite: Pi set to desktop auto-login (the service needs labwc running to have a screen);
# Restart=always keeps retrying until wayland is ready.

[Unit]
Description=Skywriter air-drawing exhibition app

[Service]
Type=simple
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-0
Environment=SDL_VIDEODRIVER=wayland
WorkingDirectory=/home/chengju
ExecStart=/home/chengju/sky/bin/python3 /home/chengju/draw_app.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
